# systemd/user/weather.service
[Unit]
Description=Update weather cache
[Service]
Type=oneshot
ExecStart=%h/.local/bin/weather-update.sh
